# deploy/jcp-mcp-server.service
[Unit]
Description=JCP MCP Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/jcp-mcp-server
ExecStart=/opt/jcp-mcp-server/jcp-mcp-server -mode=sse
Restart=always
RestartSec=5

# 环境变量配置
Environment="MCP_MODE=sse"
Environment="PORT=18080"
Environment="BASE_URL=http://193.112.101.212:18080"
Environment="LOG_LEVEL=info"

[Install]
WantedBy=multi-user.target
